# ha4linux/packaging/assets/ha4linux.service
[Unit]
Description=HA4Linux Client API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ha4linux
Group=ha4linux
EnvironmentFile=/etc/ha4linux/ha4linux.env
WorkingDirectory=/opt/ha4linux
ExecStart=/opt/ha4linux/.venv/bin/python -m app.main
Restart=on-failure
RestartSec=3
NoNewPrivileges=false
# VirtualBox VMs started via VBoxManage spawn VBoxHeadless processes that must
# survive HA4Linux service restarts and stops.
KillMode=process
# VirtualBox reports running VMs through IPC sockets under /tmp.
# PrivateTmp isolates those sockets and makes "list runningvms" return empty.
PrivateTmp=false
ProtectSystem=full
# VBoxManage requires RW access to the configured user's ~/.config/VirtualBox.
ProtectHome=false
ReadWritePaths=/var/lib/ha4linux /var/log/ha4linux /etc/ha4linux

[Install]
WantedBy=multi-user.target
