# ha4linux/packaging/assets/ha4linux.service.d/10-ha4linux-managed.conf
[Service]
# VirtualBox VMs started via VBoxManage spawn VBoxHeadless processes that must
# survive HA4Linux service restarts and stops.
KillMode=process
# VirtualBox reports running VMs through IPC sockets under /tmp.
# PrivateTmp isolates those sockets and makes "list runningvms" return empty.
PrivateTmp=false
ProtectSystem=full
# VBoxManage requires RW access to the configured user's ~/.config/VirtualBox.
ProtectHome=false
ReadWritePaths=/var/lib/ha4linux /var/log/ha4linux /etc/ha4linux
